# templates/systemd/mailman-hourly.timer.j2
[Unit]
Description=Run hourly jobs for django-admin

[Timer]
OnCalendar=hourly
Persistent=true
Unit={{ mailman_systemd_prefix }}django-admin-runjobs@hourly.service

[Install]
WantedBy=timers.target
