# templates/systemd/mailman-yearly.timer.j2
[Unit]
Description=Run hourly jobs for django-admin

[Timer]
OnCalendar=yearly
Persistent=true
Unit={{ mailman_systemd_prefix }}django-admin-runjobs@yearly.service

[Install]
WantedBy=timers.target
