# templates/systemd/mailman-weekly.timer.j2
[Unit]
Description=Run hourly jobs for django-admin

[Timer]
OnCalendar=weekly
Persistent=true
Unit={{ mailman_systemd_prefix }}django-admin-runjobs@weekly.service

[Install]
WantedBy=timers.target
