# templates/systemd/mailman-digests.service.j2
[Unit]
Description=Send periodic digests

[Service]
User={{ mailman_user }}
Group={{ mailman_group }}
Type=oneshot
ExecStart={{ mailman_path }}/bin/mailman digests --periodic
