# templates/systemd/mailman-django-admin-runjobs.service.j2
[Unit]
Description=Run django-admin jobs for mailman %I

[Service]
Type=oneshot
ExecStart={{ mailman_path }}/bin/django-admin runjobs %i
